# docker/crontab
## Ad Management System — cron schedule (all times UTC)
## Edit this file to change the schedule; rebuild container to apply.
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

## ─── Daily ───────────────────────────────────────────────────────────────────
## 06:00  Pull yesterday's performance from all active platforms
0 6 * * *   root  cd /app && python scripts/pull_performance.py >> /app/logs/cron.log 2>&1

## 06:15  Analyze + generate Change Requests (dry-run unless --apply configured)
15 6 * * *  root  cd /app && python scripts/rotate_ads.py >> /app/logs/cron.log 2>&1

## 06:30  Reconcile local spend log vs platform-reported spend
30 6 * * *  root  cd /app && python scripts/reconcile_spend.py >> /app/logs/cron.log 2>&1

## 06:45  Generate and email daily report
45 6 * * *  root  cd /app && python scripts/report.py >> /app/logs/cron.log 2>&1

## ─── Weekly ──────────────────────────────────────────────────────────────────
## Sunday 07:00  Keyword harvest: auto → exact migration, negatives
0 7 * * 0   root  cd /app && python scripts/harvest_keywords.py >> /app/logs/cron.log 2>&1

## Sunday 09:00  Cloud backup (SQLite + brand configs via rclone)
0 9 * * 0   root  cd /app && python scripts/backup.py >> /app/logs/cron.log 2>&1

## ─── (Uncomment when scripts exist) ─────────────────────────────────────────
## 0 8 * * 0   root  cd /app && python scripts/kdp_ingest.py >> /app/logs/cron.log 2>&1
